# systemd/lifelight.service
[Unit]
Description=lifelight

[Service]
ExecStart=/usr/bin/lifelight
Restart=always

[Install]
WantedBy=multi-user.target
